# oonf-dlep-radio/files/dlep_radio.init
#!/bin/sh /etc/rc.common

START=82
DAEMON='dlep_radio'

[ -n "$IPKG_INSTROOT" ] || {
        . /lib/functions/oonf_init.sh
}
